# files/pgsql_backup.cron
10 1 * * * postgres /usr/bin/pg_dumpall > /var/lib/pgsql/backups/pg_dumpall.sql && chmod 600 /var/lib/pgsql/backups/pg_dumpall.sql
